# bash/.config/profile.d/skel-override.sh
# Most of the default .profile is fine- I just need to change a few variables and it would be nice to not
# have merge conflicts every time I run skel. I don't care about /root/.profile because I should never be
# in a login shell for root anyway.
export EDITOR=nvim
export SHELL=/usr/bin/bash
export TERMINAL_COMMAND=/usr/share/sway/scripts/foot.sh
